# systemd/user/nitrogem.service
[Unit]
Description=Wallpaper server

[Service]
ExecStart=nitrogen --restore &

[Install]
WantedBy=default.target
