# systemd/user/picom.service
[Unit]
Description=Window Compositor

[Service]
ExecStart=picom

[Install]
WantedBy=default.target
